# etc/systemd/system/btrfs-snapshot-script.service
[Unit]
Description=Weekly btrfs snapshots

[Service]
Type=oneshot
ExecStart=/usr/sbin/btrfs-snapshot-script.sh